# kyros-praxis/kyros-agent.service
[Unit]
Description=Kyros Shared Agent Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=thomas
WorkingDirectory=/home/thomas/kyros-praxis/kyros-praxis
ExecStart=/usr/bin/codex --config "mcp_config=%h/mcp.json"
Restart=always
RestartSec=10
StandardOutput=file:/tmp/kyros-agent.log
StandardError=file:/tmp/kyros-agent-error.log
Environment=PATH=/usr/bin:/bin:/usr/local/bin
Environment=HOME=/home/thomas

[Install]
WantedBy=multi-user.target